# crontab
* * * * * bash /rsync_tmbackup.sh /backupSrc /backupDst &>/proc/1/fd/1
